# systemd/tidewatch.service
[Unit]
Description=Belabox Metrics Collector
Documentation=https://github.com/taniwha3/thugshells
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=metrics
Group=metrics
ExecStart=/usr/local/bin/tidewatch -config /etc/tidewatch/config.yaml

# Restart policy
Restart=always
RestartSec=10s

# Watchdog configuration (requires Type=notify integration in main.go)
# WatchdogSec=60s

# Resource limits (Milestone 2 requirements)
MemoryMax=200M
CPUQuota=20%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tidewatch

# Security hardening - Privilege restrictions
NoNewPrivileges=true
PrivateTmp=true

# Security hardening - Filesystem access
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tidewatch
ReadOnlyPaths=/etc/tidewatch

# Security hardening - Kernel features
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true

# Security hardening - Network restrictions
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true

# Security hardening - System calls
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

# Security hardening - Misc
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
PrivateDevices=false
DevicePolicy=closed
DeviceAllow=/dev/null rw
DeviceAllow=/dev/urandom r

[Install]
WantedBy=multi-user.target
